# deploy/multica-instance-router.service.example
[Unit]
Description=Multica Instance Router
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/data/multica-router/app
# Duplicate this unit per environment and point each copy to a different env file.
EnvironmentFile=/data/multica-router/multica-instance-router.env
ExecStart=/usr/local/bin/node /data/multica-router/app/server-router.mjs
Restart=always
RestartSec=5
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
